# deploy/systemd/predictor-staleness.timer
[Unit]
Description=Run predictor staleness check every 6 hours

[Timer]
OnBootSec=10min
OnUnitActiveSec=6h
Unit=predictor-staleness.service

[Install]
WantedBy=timers.target
