# deploy/systemd/predictor-daily@.service
[Unit]
Description=Predictor daily pipeline run (%i)
# One instance per domain: predictor-daily@film, predictor-daily@semiconductors,
# predictor-daily@weapons_detection.
# Notifications: run_pipeline.py pages [PREDICTOR] with a rich summary on EVERY
# outcome (success/partial/failed, incl. which stage failed) — so no generic
# OnFailure= here, which would only double-page. A hard crash that never reaches
# that ping (OOM/SIGKILL/2h-timeout) is still caught by the staleness pager
# (predictor-staleness.timer) within hours — silent no-data stalls are its job.

[Service]
Type=oneshot
User=claude
WorkingDirectory=/opt/predictor_ingest
# venv python first on PATH — the Makefile invokes bare `python`.
Environment=PATH=/opt/predictor_ingest/venv/bin:/usr/local/bin:/usr/bin:/bin
EnvironmentFile=-/opt/predictor_ingest/.env
ExecStart=/usr/bin/make daily DOMAIN=%i
TimeoutStartSec=2h
